# experiments/exp1-onechat2/flsignal.service
[Unit]
Description=Fledger signaling server
After=network.target

[Service]
Type=oneshot
ExecStart=/root/flsignal -vvv 
RemainAfterExit=yes
User=root

[Install]
WantedBy=multi-user.target
