# experiments/exp3-dummy10-router/fledger-send.service
[Unit]
Description=Fledger node
After=network.target

[Service]
Type=oneshot
ExecStart=/root/fledger --config /root/flnode -v -s ws://${FLEDGER_CENTRAL_HOST}:8765 --disable-turn-stun simulation send-chat ${FLEDGER_SEND_MSG}
RemainAfterExit=yes
User=root
EnvironmentFile=/root/env.systemd
StandardOutput=append:/var/log/fledger-send

[Install]
WantedBy=multi-user.target
