# experiments/exp1-dummy/fledger-send.service
[Unit]
Description=Fledger node
After=network.target

[Service]
Type=oneshot
ExecStart=/root/fledger --config /root/flnode -vv -s ws://10.0.0.1:8765 --disable-turn-stun --send-chat-msg message
RemainAfterExit=yes
User=root
StandardOutput=append:/var/log/fledger-send

[Install]
WantedBy=multi-user.target
