# experiments/exp3-onechat2x5/flsignal.service
[Unit]
Description=Fledger signaling server
After=network.target

[Service]
Type=simple
ExecStart=/root/flsignal -v
User=root
Restart=no
StandardOutput=append:/var/log/flsignal

[Install]
WantedBy=multi-user.target
